# automation/supervisord/programs/planner.conf
; Planner agent program configuration

[program:planner]
command=%(ENV_HOME)s/work/Depository/Depository/.venv/bin/crystalclearhouse-planner
directory=%(ENV_HOME)s/work/Depository/Depository
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/planner.err.log
stdout_logfile=/tmp/planner.out.log
user=%(ENV_USER)s
environment=PYTHONPATH="%(ENV_HOME)s/work/Depository/Depository/src"
